# supervisord.conf
[supervisord]
nodaemon=true

[program:rhodecodeinit]
command=bash -c /root/run.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
